# supervisord.conf
[supervisord]
nodaemon=true

[program:vsts]
directory=/usr/local/vsts-agent
user=vsts
command=/bin/bash -c "/usr/local/vsts-agent/start.sh"

